# systemd/codex-web.service
[Unit]
Description=Codex Web
Requires=docker.service
Wants=network-online.target
After=docker.service network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/code/codex-web
Environment=CODEX_WEB_IMAGE=codex-web:local
Environment=CODEX_WEB_CONTAINER=codex-web
Environment=CODEX_WEB_PORT=58888
Environment=CODEX_WEB_DOCKER_DATA=/root/code/codex-web/build/docker/data
Environment=CODEX_WEB_DOCKER_CODEX_HOME=/root/code/codex-web/build/docker/codex-home
Environment=CODEX_WEB_DOCKER_WORKSPACE=/root
ExecStartPre=-/usr/bin/docker rm -f ${CODEX_WEB_CONTAINER}
ExecStart=/root/code/codex-web/scripts/run-container.sh
ExecStop=/usr/bin/docker stop ${CODEX_WEB_CONTAINER}
Restart=on-failure
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
